# deploy/podman-auto-update.timer.d/override.conf
# Drop-in override for the packaged podman-auto-update.timer (which fires daily).
#
# Install once: copy to
#   ~/.config/systemd/user/podman-auto-update.timer.d/override.conf
# then `systemctl --user daemon-reload && systemctl --user enable --now podman-auto-update.timer`.
#
# The empty first OnCalendar= clears the packaged daily schedule before setting ours.

[Timer]
OnCalendar=
OnCalendar=*:0/5
RandomizedDelaySec=0
